# rootfs/etc/s6-overlay/s6-rc.d/deluge-web/run
#!/usr/bin/with-contenv sh

# load network variables
source /root/.profile

# run with limited permissions (from /usr/local/bin)
exec deluge-web \
    --loglevel info \
    --interface="$LAN_IP"
